# vagrant/zookeeper.conf
start on started networking
stop on shutdown

script
  export KAFKA_HEAP_OPTS="-Xmx192m"
  exec /opt/kafka/bin/zookeeper-server-start.sh /opt/zookeeper-ZK_PORT/config/zookeeper.properties
end script
